# deploy/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
nodaemon=true
logfile=/data/log/supervisord.log
pidfile=/data/supervisord.pid

[program:hotel_api_uwsgi]
command=uwsgi -i /code/deploy/uwsgi.ini --touch-reload=/code/deploy/wsgi.ini
stopsignal=HUP
stopasgroup=true
killasgroup=true
autorestart=true
stdout_logfile=/data/log/uwsgi.log
stderr_logfile=/data/log/uwsgi.log
stdout_logfile_maxbytes = 20MB
stderr_logfile_maxbytes = 20MB

[group:hotel_api]
programs=hotel_api_uwsgi
